# docker/php-supervisor/rootfs/etc/supervisor.d/date.ini
[program:date]
autorestart=true
command=sh -c 'sleep 1 && date && echo 666 > /dev/stderr && sleep 60'
directory=/var/www/html

redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
